# atlas-bot.service
[Unit]
Description=Atlas Service Client Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=atlas-bot
Group=atlas-bot
WorkingDirectory=/opt/atlas-bot
Environment=PATH=/opt/atlas-bot/venv/bin
ExecStart=/opt/atlas-bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-bot

# Sécurité
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/atlas-bot

[Install]
WantedBy=multi-user.target
